# roles/node_exporter/files/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart=/usr/sbin/node_exporter

[Install]
WantedBy=multi-user.target